# streamlit.service
[Unit]
Description=Streamlit Travel Planner App
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/path/to/your/app
ExecStart=/usr/local/bin/streamlit run streamlit_app.py --server.port=8501 --server.address=0.0.0.0
Restart=always
Environment="GOOGLE_API_KEY=your-api-key"

[Install]
WantedBy=multi-user.target
